# systemd/open-second-brain-enrich.service
[Unit]
Description=open-second-brain enrichment pass — LLM triplet extraction
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PATH=/home/m0nk/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/m0nk/Projects/open-second-brain
ExecStart=/home/m0nk/Projects/open-second-brain/scripts/enrich.py
StandardOutput=journal
StandardError=journal

# Security hardening: prevent OOM kill of writer (per ladybug rule #16)
OOMScoreAdjust=-200

[Install]
WantedBy=timers.target